# unified_drone.service
[Unit]
Description=Unified Drone Controller - Auto QR Detection & Payload Drop
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/addc2026
ExecStart=/usr/bin/python3 /home/pi/addc2026/unified_drone.py
Restart=always
RestartSec=5
StandardOutput=append:/home/pi/addc2026/unified_drone_stdout.log
StandardError=append:/home/pi/addc2026/unified_drone_stderr.log

# Environment variables
Environment="PYTHONUNBUFFERED=1"

# Give more time for Pixhawk to connect
StartLimitIntervalSec=200
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
